# greeks-collector.service
[Unit]
Description=Angel One Greeks Data Collector
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=zenith
Group=zenith
WorkingDirectory=/mnt/new_volume/Documents/Coding/Projects/Trading/github
Environment="PATH=/mnt/new_volume/Documents/Coding/Projects/Trading/github/venv/bin"
ExecStart=/mnt/new_volume/Documents/Coding/Projects/Trading/github/venv/bin/python greeks_collector.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
